# templates/service.ini.j2
[program:{{ service_name }}]
command={{ service_command }}
autostart = true